# etc/systemd/system/spacewx-archive.timer
[Unit]
Description=Run SpaceWX archive nightly

[Timer]
OnCalendar=*-*-* 00:05:00
Persistent=true

[Install]
WantedBy=timers.target

